# web/supervisord/mongo-connector.ini
[program:mongo-connector]
command=mongo-connector --auto-commit-interval=0 -m mongodb:27017 -t elasticsearch:9200 -d elastic2_doc_manager
numprocs=1
process_name=%(program_name)s
directory=/app
user=root
autostart=true
autorestart=true
startsecs=10
startretries=5
redirect_stderr=true
stdout_logfile=/app/mongoconnector.log
loglevel=debug
